# docker/etc/crontabs/root
* * * * * echo $(date --utc +\%Y\%m\%d_\%H\%M\%SZ) >> /var/log/php/crontest.log
